# c64-discord-bot.service
[Unit]
Description=C64 Discord Bot Service
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/path/to/c64-bot
ExecStart=/usr/bin/npm start
Restart=on-failure
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target 